# res/DEBIAN/postinst
#!/bin/bash

set -e

is_systemd() {
    [ "$(basename "$(readlink -f /proc/1/exe 2>/dev/null)" 2>/dev/null)" = "systemd" ]
}

if [ "$1" = configure ]; then
    ln -sfn /usr/share/rustdesk/rustdesk /usr/bin/rustdesk

    if is_systemd; then
        rm -f /etc/systemd/system/rustdesk.service \
              /usr/lib/systemd/system/rustdesk.service \
              /usr/lib/systemd/user/rustdesk.service
        mkdir -p /usr/lib/systemd/system/
        install -m 0644 /usr/share/rustdesk/files/systemd/rustdesk.service \
            /usr/lib/systemd/system/rustdesk.service
        systemctl daemon-reload >/dev/null 2>&1 || true
        systemctl enable rustdesk.service >/dev/null 2>&1 || true
        systemctl restart rustdesk.service >/dev/null 2>&1 || \
            systemctl start rustdesk.service >/dev/null 2>&1 || true
    fi
fi

exit 0
