# res/DEBIAN/prerm
#!/bin/bash

set -e

is_systemd() {
    [ "$(basename "$(readlink -f /proc/1/exe 2>/dev/null)" 2>/dev/null)" = "systemd" ]
}

case "$1" in
    remove|upgrade)
        rm -f /usr/bin/rustdesk

        if is_systemd; then
            if systemctl --quiet is-active rustdesk.service 2>/dev/null; then
                systemctl stop rustdesk.service >/dev/null 2>&1 || true
            fi
            if systemctl --quiet is-enabled rustdesk.service 2>/dev/null; then
                systemctl disable rustdesk.service >/dev/null 2>&1 || true
            fi
            rm -f /etc/systemd/system/rustdesk.service \
                  /usr/lib/systemd/system/rustdesk.service \
                  /usr/lib/systemd/user/rustdesk.service

            serverUser=$(ps -ef | grep -E 'rustdesk +--server' | grep -v 'sudo ' | awk '{print $1}' | head -1)
            if [ -n "$serverUser" ] && [ "$serverUser" != "root" ]; then
                systemctl --machine=${serverUser}@.host --user stop rustdesk >/dev/null 2>&1 || true
            fi
            systemctl daemon-reload >/dev/null 2>&1 || true
        fi
        ;;
esac

exit 0
